# home_automation/WifeSignal/ops/com.ariemeir.wife-signal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ariemeir.wife-signal</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ariemeir/dev/scripts/home_automation/WifeSignal/server/service.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ariemeir/dev/scripts/home_automation/WifeSignal/server</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/ariemeir/.wife-signal/logs/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ariemeir/.wife-signal/logs/server-error.log</string>
</dict>
</plist>
